# etc/cron.d/backupninja.in
# /etc/cron.d/backupninja -- cron tab entry for package backupninja

PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# run backupninja every hour on the hour
0 * * * * root if [ -x @sbindir@/backupninja ]; then @sbindir@/backupninja; fi
