# etc/cron.d/backupninja
# cron job for backupninja (every hour, on the hour)
00 * * * * root if [ -x /usr/sbin/backupninja -a -f /etc/backupninja.conf ]; then /usr/sbin/backupninja; fi
