# etc/cron.d/backupninja
# cron job for backupninja (once a day)
00 01 * * * root if [ -x /usr/sbin/backupninja -a -f /etc/backupninja.conf ]; then /usr/sbin/backupninja; fi
